# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server

[program:nginx]
command=/usr/local/openresty/nginx/sbin/nginx